# s6-overlay/s6-rc.d/agent-bootstrap/up
#!/bin/bash
set -euo pipefail

# Fix ownership of the VS Code server directory tree.
# When Docker mounts the agent-home named volume and THEN applies the
# ./vscode-server/extensions bind-mount, if .vscode-server does not yet
# exist in the volume Docker creates it as root.  That makes the directory
# unwritable by the agent SSH user, causing a "Permission denied" error
# when VS Code Remote-SSH tries to SCP its server binary into
# /home/agent/.vscode-server/.
#
# We must NOT chown .vscode-server/extensions because that path is a
# bind-mount from the host – touching its ownership would change the host
# directory permissions.  Likewise for .vscode-server/data/Machine/settings.json.
VSCODE_DIRS=(
    /home/agent/.vscode-server
    /home/agent/.vscode-server/data
    /home/agent/.vscode-server/data/Machine
)

mkdir -p "${VSCODE_DIRS[@]}"
chown agent:agent "${VSCODE_DIRS[@]}"

exec /command/with-contenv s6-setuidgid agent \
    /usr/local/share/copilot-code-server-container/agent-bootstrap.sh true
